# example.service
[Unit]
Description=Example service

[Service]
ExecStart=/root/fake.sh

[Install]
WantedBy=multi-user.target
Alias=example.service

