# fake.service
[Unit]
Description=Fake service

[Service]
ExecStart=/root/fake.sh

[Install]
WantedBy=multi-user.target
Alias=fake.service

