# paged.service
[Unit]
Description=Pagekite service

[Service]
ExecStart=/root/page.sh


[Install]
WantedBy=multi-user.target
Alias=5555.service

